# sandbox/image/overlay/etc/init.d/agentd
#!/sbin/openrc-run

name="agentd"
description="Cowork sandbox agent"

# Use fixed paths to avoid variable expansion issues at script parse time
command=/usr/bin/node
command_args="/opt/agent-runner/index.js"
command_background=yes
pidfile=/run/agentd.pid
directory=/workspace

# Log output for debugging
output_log=/var/log/agentd.log
error_log=/var/log/agentd.log

depend() {
    # Start as early as possible — only need /dev populated and kernel modules loaded.
    # Networking is started concurrently; agent-runner handles network unavailability.
    after modules
    use net
}

start_pre() {
    # Ensure directories exist
    checkpath -d -m 0755 /workspace
    checkpath -d -m 0755 /workspace/ipc
    checkpath -f -m 0644 /var/log/agentd.log

    # Load 9p kernel modules
    modprobe 9p 2>/dev/null || true
    modprobe 9pnet 2>/dev/null || true
    modprobe 9pnet_virtio 2>/dev/null || true
    # Load virtio-serial module (used as IPC fallback on Windows)
    modprobe virtio_console 2>/dev/null || true

    # Try to mount IPC share directory via 9p
    # If mount fails, agent-runner will fall back to virtio-serial IPC
    if ! mountpoint -q /workspace/ipc; then
        if mount -t 9p -o trans=virtio,version=9p2000.L,msize=65536 ipc /workspace/ipc 2>/dev/null; then
            einfo "IPC share mounted via 9p"
        else
            ewarn "9p mount unavailable, will use virtio-serial IPC"
        fi
    fi
}
